# deploy/dayos-sync.service
[Unit]
Description=Sync the agent's memory banks (DayOS + playbook)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=agent
WorkingDirectory=/opt/instatank-agent
EnvironmentFile=/opt/instatank-agent/.env
# '-' prefix: a playbook failure is recorded in its own status file (the bot
# warns loudly) but must never block the DayOS sync below it.
ExecStart=-/opt/instatank-agent/venv/bin/python playbook_sync.py
ExecStart=/opt/instatank-agent/venv/bin/python dayos_sync.py

# Same hardening profile as the bot: may only write under memory/
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/instatank-agent/memory
ProtectHome=yes
PrivateTmp=yes
