# deploy/memory-backup.timer
[Unit]
Description=Nightly brain-memory backup (03:30 IST)

[Timer]
OnCalendar=*-*-* 03:30:00 Asia/Kolkata
Persistent=true

[Install]
WantedBy=timers.target
